# roles/mailserver/dkim/files/rc/opendkim
#!/bin/sh

# PROVIDE: opendkim
# REQUIRE: DAEMON
# BEFORE: mail
# KEYWORD: shutdown

#   -------------------------------------------------------------
#   OpenDKIM
#   - - - - - - - - - - - - - - - - - - - - - - - - - - - - - - -
#   Project:        Nasqueron
#   License:        Trivial work, not eligible to copyright
#   Source file:    roles/mailserver/dkim/files/rc/opendkim
#   -------------------------------------------------------------
#
#   <auto-generated>
#       This file is managed by our rOPS SaltStack repository.
#
#       Changes to this file may cause incorrect behavior
#       and will be lost if the state is redeployed.
#   </auto-generated>

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
# opendkim_enable (bool): Set it to YES to enable the service.
#           Default is "NO"
# opendkim_user (user):   Set user to run opendkim
#           Default is "opendkim".
# opendkim_group (group):   Set user to run opendkim
#           Default is "mail".
# opendkim_config (config file path):   Set config to run opendkim
#           Default is "/usr/local/etc/opendkim/opendkim.conf".
# opendkim_perms (perms):   Set config to run opendkim
#           Default is "0755".

. /etc/rc.subr

name="opendkim"
rcvar="${name}_enable"

load_rc_config $name

: ${opendkim_enable:="NO"}
: ${opendkim_user:="opendkim"}
: ${opendkim_group:="mail"}
: ${opendkim_config:="/usr/local/etc/opendkim/opendkim.conf"}
: ${opendkim_perms:="0755"}

pidfile="/var/run/opendkim/opendkim.pid"
socketfile="/var/run/opendkim/opendkim.sock"
command="/usr/local/sbin/opendkim"
command_args="-l -u ${opendkim_user}:${opendkim_group} -P ${pidfile} -x ${opendkim_config}"

start_precmd=opendkim_exec_startprecmd
stop_postcmd=opendkim_exec_stop_postcmd

opendkim_exec_startprecmd()
{
    mkdir -p /var/run/opendkim
    chmod ${opendkim_perms} /var/run/opendkim
    chown ${opendkim_user}:${opendkim_group} /var/run/opendkim
}

opendkim_exec_stop_postcmd()
{
    rm -f ${pidfile}
    rm -f ${socketfile}
}

run_rc_command "$1"
