# roles/mailserver/systemd-unit/files/lxc-container-mailserver.service
#   -------------------------------------------------------------
#   Nasqueron mail services
#   - - - - - - - - - - - - - - - - - - - - - - - - - - - - - - -
#   Project:        Nasqueron
#   Created:        2015-12-16
#   License:        Trivial work, not eligible to copyright
#   Source file:    roles/mailserver/systemd-unit/files/lxc-container-mailserver.service
#   -------------------------------------------------------------
#
#   <auto-generated>
#       This file is managed by our rOPS SaltStack repository.
#
#       Changes to this file may cause incorrect behavior
#       and will be lost if the state is redeployed.
#   </auto-generated>

[Unit]
Description=LXC mail server container

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/system/lxc-container-mailserver-start
ExecStop=/usr/lib/systemd/system/lxc-container-mailserver-stop

[Install]
WantedBy=multi-user.target
