# roles/core/letsencrypt-renew/files/letsencrypt.service
[Unit]
Description=Renew Let's encrypt certificates.

[Service]
Type=simple
ExecStart=letsencrypt-renewal
