# roles/paas-docker/systemd-unit/files/docker-containers.service
[Service]
Type=simple
RemainAfterExit=yes
ExecStart=docker_start
ExecStop=docker_stop
