# roles/paas-docker/systemd-unit/files/letsencrypt.timer
[Unit]
Description=Runs letsencrypt-renewal every month

[Timer]
OnCalendar=*-*-26 12:15:00
Persistent=yes
Unit=letsencrypt.service
